# scripts/systemd/llm-vault-sync.service.example
# scripts/systemd/llm-vault-sync.service.example
# Linux(systemd user) 자동 동기화 서비스. timer가 주기적으로 트리거한다.
#
# 설치(기기마다 1회):
#   mkdir -p ~/.config/systemd/user
#   cp scripts/systemd/llm-vault-sync.service.example ~/.config/systemd/user/llm-vault-sync.service
#   cp scripts/systemd/llm-vault-sync.timer.example   ~/.config/systemd/user/llm-vault-sync.timer
#   # 두 파일의 <REPO> 를 이 기기의 저장소 절대경로로 치환
#   systemctl --user daemon-reload
#   systemctl --user enable --now llm-vault-sync.timer
#   loginctl enable-linger "$USER"        # 로그아웃 상태에서도 타이머 동작
#   systemctl --user list-timers | grep llm-vault     # 등록 확인
#   journalctl --user -u llm-vault-sync.service -f     # 로그
# 제거:
#   systemctl --user disable --now llm-vault-sync.timer
#
# 자격증명: HTTPS면 git credential helper(libsecret/store) 또는 PAT, SSH면 ssh-agent
#           등 비대화식 인증이 설정돼 있어야 한다(백그라운드 실행이라 프롬프트 불가).

[Unit]
Description=llm-vault git auto-sync (oneshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=<REPO>
ExecStart=/bin/bash <REPO>/scripts/sync.sh
