# scripts/systemd/llm-vault-sync.timer.example
# scripts/systemd/llm-vault-sync.timer.example
# Linux(systemd user) 자동 동기화 타이머. 15분 간격 + 부팅 후 2분.
# 설치 절차는 llm-vault-sync.service.example 헤더 참조.

[Unit]
Description=Run llm-vault git auto-sync every 15 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true
Unit=llm-vault-sync.service

[Install]
WantedBy=timers.target
